# assets/bootstrap/systemd/hytale-auth-scan.service
[Unit]
Description=Scan Hytale logs for auth URLs and notify Discord
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=hytale
Group=hytale
WorkingDirectory=/opt/hytale
ExecStart=/opt/hytale/bin/hytale-auth-scan.sh

