# assets/bootstrap/systemd/hytale-backup-sync.timer
[Unit]
Description=Periodic Hytale backup upload to S3

[Timer]
OnBootSec=10min
OnUnitActiveSec=30min
Persistent=true

[Install]
WantedBy=timers.target

