# assets/bootstrap/systemd/hytale-backup-sync.service
[Unit]
Description=Sync Hytale backups to S3
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/hytale/bin/hytale-backup-sync.sh

